# ubuntu22.04/autorun-clash.service
[Unit]
Description=My Program

[Service]
ExecStart=/etc/systemd/system/autorun.sh
Restart=always

[Install]
WantedBy=multi-user.target
